# root/etc/profile.d/00-aliases.sh
alias ll="LC_COLLATE=C ls -ahl --group-directories-first --color=auto"
alias ..="cd .."